# templates/supervisord.service.erb
# puppet managed file. local changes will be overwritten.
[Unit]
Description=Supervisor process control system for UNIX
Documentation=http://supervisord.org
After=network.target<% if !@_after.empty? %> <%= @_after.join(" ") %><% end %>

[Service]
Type=<%= @_type %>
ExecStart=<%= @_binpath %> -c <%= @_confdir %>/<%= @_conffile %>
ExecStop=<%= @_ctlpath %> $OPTIONS shutdown
ExecReload=<%= @_ctlpath %> -c <%= @_confdir %>/<%= @_conffile %> $OPTIONS reload
KillMode=<%= @_killmode %>
Restart=<%= @_restart %>
RestartSec=<%= @_restartsec %>

[Install]
WantedBy=multi-user.target
